# scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24
autostart=true
autorestart=true
priority=5
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err

[program:vnc]
command=/usr/bin/vncserver :1 -geometry 1920x1080 -depth 24
autostart=true
autorestart=true
priority=10
environment=DISPLAY=":1"

[program:novnc]
command=websockify --web /usr/share/novnc 6080 localhost:5901
autostart=true
autorestart=true
priority=15

[program:ccb]
command=/bin/bash -c "while true; do ccb || sleep 5; done"
autostart=true
autorestart=true
priority=20
